# systemd-units/kf2server-update.timer
[Unit]
Description=Update the Killing Floor 2 Server

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
